# packages/linux-drivers/RTL8723BS_BT/system.d/rtk_bt.timer
[Unit]
Description=Realtek RTL8723BS bluetooth service startup after delay

[Timer]
OnBootSec=10

[Install]
WantedBy=graphical.target
